# support/systemd/nctu5gc-upfd.service.in
[Unit]
Description=NCTU5GC UPF Daemon
After=networking.service

[Service]
Type=simple

User=nctu5gc
Group=nctu5gc

RuntimeDirectory=nctu5gc-upfd
PIDFile=@LIB_DIR@/run/nctu5gc-upfd/pid

Restart=always
ExecStart=@BIN_DIR@/nctu5gc-upfd -f @SYSCONF_DIR@/nctu5gc/upf.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
