# support/systemd/nctu5gc-amfd.service.in
[Unit]
Description=NCTU5GC AMF Daemon
After=networking.service

[Service]
Type=simple

User=nctu5gc
Group=nctu5gc

RuntimeDirectory=nctu5gc-amfd
PIDFile=@LIB_DIR@/run/nctu5gc-amfd/pid

Restart=always
ExecStart=@BIN_DIR@/nctu5gc-amfd -f @SYSCONF_DIR@/nctu5gc/amf.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
