# support/systemd/nextepc-smfd.service.in
[Unit]
Description=NextEPC SMF Daemon
After=networking.service

[Service]
Type=simple

User=nextepc
Group=nextepc

RuntimeDirectory=nextepc-smfd
PIDFile=@LIB_DIR@/run/nextepc-smfd/pid

Restart=always
ExecStart=@BIN_DIR@/nextepc-smfd -f @SYSCONF_DIR@/nextepc/smf.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
