# support/systemd/nextepc-upfd.service.in
[Unit]
Description=NextEPC P-GWU Daemon
After=networking.service

[Service]
Type=simple

User=nextepc
Group=nextepc

RuntimeDirectory=nextepc-upfd
PIDFile=@LIB_DIR@/run/nextepc-upfd/pid

Restart=always
ExecStart=@BIN_DIR@/nextepc-upfd -f @SYSCONF_DIR@/nextepc/upf.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
