# support/systemd/nextepc-hssd.service.in
[Unit]
Description=NextEPC HSS Daemon
BindTo=mongodb.service
After=networking.service mongodb.service

[Service]
Type=simple

User=free5gc
Group=free5gc

RuntimeDirectory=nextepc-hssd
PIDFile=@LIB_DIR@/run/nextepc-hssd/pid

Restart=always
ExecStart=@BIN_DIR@/nextepc-hssd -f @SYSCONF_DIR@/free5gc/hss.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
