# support/systemd/open5gc-smfd.service.in
[Unit]
Description=NextEPC SMF Daemon
After=networking.service

[Service]
Type=simple

User=nextepc
Group=nextepc

RuntimeDirectory=open5gc-smfd
PIDFile=@LIB_DIR@/run/open5gc-smfd/pid

Restart=always
ExecStart=@BIN_DIR@/open5gc-smfd -f @SYSCONF_DIR@/open5gc/smf.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
