# support/systemd/open5gc-pcrfd.service.in
[Unit]
Description=NextEPC PCRF Daemon
BindTo=mongodb.service
After=networking.service mongodb.service

[Service]
Type=simple

User=nextepc
Group=nextepc

RuntimeDirectory=open5gc-pcrfd
PIDFile=@LIB_DIR@/run/open5gc-pcrfd/pid

Restart=always
ExecStart=@BIN_DIR@/open5gc-pcrfd -f @SYSCONF_DIR@/open5gc/pcrf.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
